# teamspeak3-server/daemon
#! /bin/sh

# Original Author: Bruno Pereira
# Modified by: Maciej Krüger

### BEGIN INIT INFO
# Provides:          teamspeak3
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts, stops implements a fake restart and a database reset for teamspeak3 server
### END INIT INFO

dir=$SNAP/bin/server
data=$SNAP_USER_DATA/server
mkdir -p $data
pidfile=$data/ts3server.pid

doscript() {
	if [ -e $pidfile ]; then
		if ( kill -0 $(cat $pidfile) 2> /dev/null ); then
			run="true"
		else
			run="false"
			rm $pidfile
		fi
	else
		run="false"
	fi
	case "$1" in
		start)
			if [ $run == "true" ]; then
				echo "The server is already running, try restart or stop"
				exit 1
			fi
			linkall
			#Starts teamspeak3 server and send all messages to /dev/null.
			echo "Starting teamspeak3 server"
			$data/ts3server_startscript.sh start > /dev/null
			;;
		stop)
			#Shuts down teamspeak3 server
			echo "Stopping teamspeak3 server"
			if [ -f $data/ts3server_startscript.sh ]; then
				sh $data/ts3server_startscript.sh stop > /dev/null
			fi
			sleep 0.5s
			unlinkall
			;;
		restart)
			#Server restart
			echo "Restarting teamspeak3 server"
			if [ $run == "true" ]; then
				doscript stop
			fi
			doscript start
			;;
		init_db|reset)
			#Removes teamspeak3 server database and creates a new one with defaults
			LOCAL_DB_FILE="$data/ts3server.sqlitedb"
			LOCAL_PID_FILE="$data/ts3server.pid"
			#If server is running, stop it
			if [ -f "$LOCAL_PID_FILE" ]; then
				sh $data/ts3server_startscript.sh stop > /dev/null 2> /dev/null
				sleep 5s
			fi
			linkall
			#Delete the current database if present
			if [ -f "$LOCAL_DB_FILE" ]; then
				rm $data/ts3server.sqlitedb
			fi
			sh $data/ts3server_startscript.sh start > /dev/null 2> /dev/null
			sleep 5s
			sh $data/ts3server_startscript.sh stop > /dev/null 2> /dev/null

			pw=$(cat /dev/urandom | tr -dc 'a-zA-Z0-9' | fold -w 10 | head -n 1)
			token=$(cat /dev/urandom | tr -dc 'a-zA-Z0-9' | fold -w 32 | head -n 1)

			sqlite3 $data/ts3server.sqlitedb "UPDATE tokens SET token_key='$token' WHERE server_id='1'"
			sqlite3 $data/ts3server.sqlitedb "UPDATE server_properties SET value='$token' WHERE ident='virtualserver_autogenerated_privilegekey'"
			sqlite3 $data/ts3server.sqlitedb "UPDATE clients SET client_login_password='$pw' WHERE client_login_name='serveradmin'"
			echo "
------------------------------------------------------------------
                      I M P O R T A N T
------------------------------------------------------------------
               Server Query Admin Account created
         Username \"serveradmin\" Password \"$pw\"
------------------------------------------------------------------
                    Server 1 Admin Token
             $token
------------------------------------------------------------------
	" > $data/READ_ME_IMPORTANT.txt
			;;
		*)
			echo "Usage: $0 {start|stop|restart|init_db}"
			exit 1
	esac
}

files="$dir/*.so $dir/redist $dir/sql $dir/ts3server* $dir/VERSION"

unlinkall() {
	for f in $files; do
		rm -f ${f/$dir/$data}
	done
}

linkall() {
	unlinkall
	for f in $files; do
		ln -s $f ${f/$dir/$data}
	done
	rm -f $data/ts3server*.sh
	cp $dir/ts3server*.sh $data
}

if ! [ -f $data/ts3server.sqlitedb ]; then
	echo -n "Setting up teamspeak3 server"
	doscript init_db
	echo " done"
fi

doscript $1

exit 0
