# teamspeak3-server/daemon
#! /bin/sh

# Original Author: Bruno Pereira
# Modified by: Maciej Krüger

### BEGIN INIT INFO
# Provides:          teamspeak3
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts, stops implements a fake restart and a database reset for teamspeak3 server
### END INIT INFO

dir=$SNAP/bin/server
data=$SNAP_USER_DATA/server
mkdir -p $data

doscript() {
	case "$1" in
	start)
		#Starts teamspeak3 server and send all messages to /dev/null.
		echo "Starting teamspeak3 server"
		#mount -t tmpfs tmpfs /dev/shm
		$dir/ts3server_startscript.sh start > /dev/null 2> /dev/null
		;;
	stop)
		#Shuts down teamspeak3 server
		echo "Stopping teamspeak3 server"
		sh $dir/ts3server_startscript.sh stop > /dev/null 2> /dev/null
		sleep 5s
		#umount /dev/shm
		;;
	restart)
		#Server restart
		echo "Restarting teamspeak3 server"
		sh $dir/ts3server_startscript.sh restart > /dev/null 2> /dev/null
		;;
	init_db|reset)
		#Removes teamspeak3 server database and creates a new one with defaults
		LOCAL_DB_FILE="$dir/ts3server.sqlitedb"
		LOCAL_PID_FILE="$dir/ts3server.pid"
		#If server is running, stop it
		if [ -f "$LOCAL_PID_FILE" ]; then
			sh $dir/ts3server_startscript.sh stop > /dev/null 2> /dev/null
			sleep 5s
		fi
		#Delete the current database if present
		if [ -f "$LOCAL_DB_FILE" ]; then
			rm $dir/ts3server.sqlitedb
		fi
		sh $dir/ts3server_startscript.sh start > /dev/null 2> /dev/null
		sleep 10s
		sh $dir/ts3server_startscript.sh stop > /dev/null 2> /dev/null

		pw=$(cat /dev/urandom | tr -dc 'a-zA-Z0-9' | fold -w 10 | head -n 1)
		token=$(cat /dev/urandom | tr -dc 'a-zA-Z0-9' | fold -w 32 | head -n 1)

		sqlite3 $data/ts3server.sqlitedb "UPDATE tokens SET token_key='$token' WHERE server_id='1'"
		sqlite3 $data/ts3server.sqlitedb "UPDATE server_properties SET value='$token' WHERE ident='virtualserver_autogenerated_privilegekey'"
		sqlite3 $data/ts3server.sqlitedb "UPDATE clients SET client_login_password='$pw' WHERE client_login_name='serveradmin'"
		echo "
		============[IMPORTANT]============
		Query Username: serveradmin
		Query Password: $pw

		Server 1 Admin Token: $token

		Keep this file
			and its contents secret!!!
		===================================
		" > $data/READ_ME_IMPORTANT.txt
		;;
	*)
		echo "Usage: $0 {start|stop|restart|init_db}"
		exit 1
	esac
}

linkall() {
	for f in $dir/*.so $dir/redist $dir/sql $dir/ts3server* $dir/VERSION; do
		rm -f ${f/$dir/$data}
		ln -s $f ${f/$dir/$data}
	done
	rm -f $data/ts3server*.sh
	cp $dir/ts3server*.sh $data
}

linkall

dir=$data

if ! [ -f $data/ts3server.sqlitedb ]; then
	echo -n "Setting up teamspeak3 server"
	doscript init_db
	echo " done"
fi

doscript $1

exit 0
